# contrib/init/rezabitd.conf
description "Rezabit Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env REZABITD_BIN="/usr/bin/rezabitd"
env REZABITD_USER="rezabit"
env REZABITD_GROUP="rezabit"
env REZABITD_PIDDIR="/var/run/rezabitd"
# upstart can't handle variables constructed with other variables
env REZABITD_PIDFILE="/var/run/rezabitd/rezabitd.pid"
env REZABITD_CONFIGFILE="/etc/rezabit/rezabit.conf"
env REZABITD_DATADIR="/var/lib/rezabitd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # rezabitd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$REZABITD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run rezabitd."
        echo "The setting must appear in $REZABITD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Rezabit Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$REZABITD_PIDDIR"
    chmod 0755 "$REZABITD_PIDDIR"
    chown $REZABITD_USER:$REZABITD_GROUP "$REZABITD_PIDDIR"
    chown $REZABITD_USER:$REZABITD_GROUP "$REZABITD_CONFIGFILE"
    chmod 0660 "$REZABITD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$REZABITD_PIDFILE" \
    --chuid $REZABITD_USER:$REZABITD_GROUP \
    --exec "$REZABITD_BIN" \
    -- \
    -pid="$REZABITD_PIDFILE" \
    -conf="$REZABITD_CONFIGFILE" \
    -datadir="$REZABITD_DATADIR" \
    -disablewallet \
    -daemon

